# deploy/biggest-losers-api.service
[Unit]
Description=Biggest Losers API (FastAPI + scheduler)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{ROOT}}/backend
EnvironmentFile={{ROOT}}/backend/.env
ExecStart={{ROOT}}/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
